# configs/supervisor/cron.conf
[program:cron.super]
user = root
username = root
command = cron -f -L 15
autostart=true
autorestart=true

stderr_logfile=/var/log/cron.errors.log
stdout_logfile=/var/log/cron.console.log
